# sejmbot.service
[Unit]
Description=SejmBot - Parser transkryptów Sejmu RP (Daemon)
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=3

[Service]
Type=simple
Restart=on-failure
RestartSec=60
User=filip
Group=filip
WorkingDirectory=/home/filip/sejmbot
Environment=PATH=/home/filip/sejmbot/venv/bin
# Uruchamiaj główny bot jako daemon (nie scheduler!)
ExecStart=/home/filip/sejmbot/venv/bin/python sejmbot.py --daemon
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

# Logi
StandardOutput=append:/home/filip/sejmbot/logs/daemon.log
StandardError=append:/home/filip/sejmbot/logs/daemon.log

# Ograniczenia zasobów dla Pi
CPUQuota=90%
MemoryMax=512M
OOMScoreAdjust=100

# Bezpieczeństwo
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target